# aws/apitools/cfn-init-2.0-20/init/systemd/cfn-hup.service
[Unit]
Description=cfn-hup daemon

[Service]
ExecStart=/usr/bin/cfn-hup -v
PIDFile=/var/run/cfn-hup.pid

[Install]
WantedBy=multi-user.target